# systemd/announcements-display.timer
# /etc/systemd/system/announcements-display.timer
#
# Runs announcements-display.service every minute.
# The service executes schedule_display.sh, which performs a long-running loop.

[Unit]
Description=Announcements: schedule display power every minute

[Timer]
OnBootSec=60
OnUnitActiveSec=60
Unit=announcements-display.service
Persistent=true

[Install]
WantedBy=timers.target

