# systemd/announcements-display.service
# /etc/systemd/system/announcements-display.service
#
# Runs announcements-display.sh
# The script itself loops forever and handles display ON/OFF polling.

[Unit]
Description=Announcements: display power scheduler
After=network.target

[Service]
Type=simple
ExecStart=/srv/announcements/announcements-display.sh
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
