# systemd/announcements-status.service
# /etc/systemd/system/announcements-status.service
#
# Runs announcements-status.sh continuously to update timestamped
# status files (_STATUS_<timestamp>.txt) in the inbox while a
# conversion is in progress. The script contains its own loop
# (sleep 30 between updates).
#
# Must run as root so it can always read the temp/staging.*
# directories, regardless of ownership/mode.

[Unit]
Description=Announcements: status watcher
After=network.target

[Service]
Type=simple
ExecStart=/srv/announcements/announcements-status.sh
Restart=always
RestartSec=5
User=root
Group=root

[Install]
WantedBy=multi-user.target